# claude-code/rootfs/etc/s6-overlay/s6-rc.d/claude-code/run
#!/command/with-contenv bash
# shellcheck shell=bash

export HOME=/data
export LANG=C.UTF-8
export LC_ALL=C.UTF-8

# ── Create persistent workspace ──────────────────────────
mkdir -p /data/workspace
mkdir -p /data/.claude

# ── Write CLAUDE.md ──────────────────────────────────────
cat > /data/workspace/CLAUDE.md << 'CLAUDEMD'
# Home Assistant Claude Code

You are running inside a Home Assistant add-on on a Raspberry Pi.

## HA API Access (MCP Tools)

- `ha_get_states` / `ha_get_state` — Read entities
- `ha_call_service` — Control devices (ALWAYS use entity_id, not area_id)
- `ha_get_services` — List services
- `ha_create_automation` — Create automations
- `ha_get_config` — HA config
- `ha_fire_event` — Fire events
- `ha_addon_info` / `ha_supervisor_info` — System info

## Dashboard Modifications (WebSocket)
```python
import asyncio, json, os, websockets
async def main():
    token = os.environ['SUPERVISOR_TOKEN']
    async with websockets.connect('ws://supervisor/core/websocket') as ws:
        await ws.recv()
        await ws.send(json.dumps({'type': 'auth', 'access_token': token}))
        await ws.recv()
```

## Environment
- Working dir: /data/workspace (persistent, backed up)
- HA config: /config (read-only)
- SUPERVISOR_TOKEN available in env
CLAUDEMD

# ── Write MCP config ─────────────────────────────────────
cat > /data/workspace/.mcp.json << MCPJSON
{
  "mcpServers": {
    "home-assistant": {
      "command": "/root/.bun/bin/bun",
      "args": ["run", "/app/src/ha-mcp.ts"],
      "env": {
        "SUPERVISOR_TOKEN": "${SUPERVISOR_TOKEN}"
      }
    }
  }
}
MCPJSON

# ── Apply credentials ────────────────────────────────────
if [ -f /data/options.json ]; then
    CREDS_JSON=$(jq -r '.credentials_json // ""' /data/options.json)
    if [ -n "$CREDS_JSON" ] && [ "$CREDS_JSON" != "" ] && [ "$CREDS_JSON" != "null" ]; then
        echo "$CREDS_JSON" > /data/.claude/.credentials.json
        echo "Credentials applied."
    fi
fi

# ── Start ttyd on internal port with base-path ───────────
echo "Starting ttyd on port 5101..."
ttyd \
    --port 5101 \
    --base-path /terminal \
    --writable \
    /usr/local/bin/claude-launcher.sh &

sleep 1

# ── Start Bun web server on Ingress port ─────────────────
export CLAUDE_WORKING_DIR="/data/workspace"
cd /app
echo "Starting web server on port 5100..."
exec /root/.bun/bin/bun run src/server.ts
